# my-skills/my-skills-sync/my-skills-sync.service
[Unit]
Description=My Skills Sync Service - Sync my-* skills to local directory
After=network.target

[Service]
Type=forking
User=root
WorkingDirectory=/root/sh/my-skills/sync-my-skills
ExecStart=/root/sh/my-skills/sync-my-skills/watch-skills.sh start
ExecStop=/root/sh/my-skills/sync-my-skills/watch-skills.sh stop
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
